# productivity/toggle-leisure.service
[Unit]
Description=toggle leisure service

[Service]
ExecStart=$HOME/pet/dotfiles/productivity/toggle-leisure.sh
