# productivity/toggle-leisure.timer
[Unit]
Description=toggle leisure timer

[Timer]
OnCalendar=Mon..Sat 17:00
Persistent=true

[Install]
WantedBy=timers.target
