# productivity/toggle-leisure.service
[Unit]
Description=toggle leisure service

[Service]
ExecStart=/home/ayorgo/pet/dotfiles/productivity/toggle-leisure.sh
