# productivity/toggle-leisure.timer
[Unit]
Description=toggle leisure timer

[Timer]
OnCalendar=Mon..Sat 14:30
Persistent=true

[Install]
WantedBy=timers.target
