# productivity/toggle-work.service
[Unit]
Description=toggle work service

[Service]
ExecStart=/home/ayorgo/pet/dotfiles/productivity/toggle-work.sh
