# productivity/toggle-work.service
[Unit]
Description=toggle work service

[Service]
ExecStart=$HOME/pet/dotfiles/productivity/toggle-work.sh
